# crontab
# Restart logmanager
0 0 * * *	root	/usr/bin/supervisorctl restart logmanager
