# supervisord.conf
[supervisord]
nodaemon=true

[program:emobility]
directory=/data/eMobility/
command=/home/steve/eMobility/startEMobility.sh
autorestart=true


