# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
user=mysql
command=/usr/sbin/mysqld
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:centralsystem]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 10 && ifconfig eth0:0 10.0.1.5 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/centralsystem/target/CentralSystem.jar'
autorestart=true

[program:chargepoint1]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 100 && ifconfig eth0:1 10.0.1.10 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_1.xml -b true -s true'
autorestart=true

[program:chargepoint2]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 105 && ifconfig eth0:2 10.0.1.11 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_2.xml -b true -s true'
autorestart=true

[program:chargepoint3]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 110 && ifconfig eth0:3 10.0.1.12 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_3.xml -b true -s true'
autorestart=true

[program:chargepoint4]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 115 && ifconfig eth0:4 10.0.1.13 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_4.xml -b true -s true'
autorestart=true

[program:chargepoint5]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 120 && ifconfig eth0:5 10.0.1.14 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_5.xml -b true -s true'
autorestart=true

[program:chargepoint6]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 125 && ifconfig eth0:6 10.0.1.15 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_6.xml -b true -s true'
autorestart=true

[program:chargepoint7]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 130 && ifconfig eth0:7 10.0.1.16 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_7.xml -b true -s true'
autorestart=true

[program:chargepoint8]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 135 && ifconfig eth0:8 10.0.1.17 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_8.xml -b true -s true'
autorestart=true

[program:chargepoint9]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 140 && ifconfig eth0:9 10.0.1.18 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_9.xml -b true -s true'
autorestart=true

[program:chargepoint10]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 145 && ifconfig eth0:10 10.0.1.19 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_10.xml -b true -s true'
autorestart=true

[program:chargepoint11]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 150 && ifconfig eth0:11 10.0.1.20 netmask 255.255.255.0 && exec java -jar /opt/emobility/src/chargepoint/target/ChargePoint.jar -c /opt/emobility/conf/chargepoint_11.xml -b true -s true'
autorestart=true

[program:logmanager]
directory=/opt/emobility/
command=/bin/bash -c 'sleep 155 && /opt/emobility/scripts/startLog.sh'
autorestart=true

[program:ewsposter]
command=/usr/bin/python /opt/ewsposter/ews.py -c /data/ews/conf/ -j /data/ews/emobility/ -m emobility -l 60
autorestart=true
