# scripts/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.
#Restart tcpdump at 12:00 AM (kill it and supervisord will restart it)
at midnight <<< 'pkill tcpdump'
exit 0
